# backend/systemd/atlas-pipeline.timer
[Unit]
Description=ATLAS Pipeline Timer — run daily at 19:00 IST (13:30 UTC)
Requires=atlas-pipeline.service
Documentation=https://github.com/jslwealth/atlas

[Timer]
# 19:00 IST = 13:30 UTC (IST is UTC+5:30)
OnCalendar=*-*-* 13:30:00 UTC
AccuracySec=60
Persistent=true
Unit=atlas-pipeline.service

[Install]
WantedBy=timers.target
